# Linux/udev/98-action.rules
SUBSYSTEM=="block", ACTION=="add", KERNEL=="sd[a-z]", RUN+="/tmp/action.sh"